# etc/cron.monthly/standard
#!/bin/sh
# /etc/cron.monthly/standard: standard monthly maintenance script

# rotation of wtmp and btmp taken over by logrotate
