# var/lib/dpkg/info/apache.list
/.
/etc
/etc/apache
/etc/init.d
/etc/init.d/apache
/etc/cron.daily
/etc/cron.daily/apache
/usr
/usr/sbin
/usr/sbin/apache
/usr/sbin/apacheconfig
/usr/sbin/apachectl
/usr/lib
/usr/lib/cgi-bin
/usr/share
/usr/share/doc
/usr/share/doc/apache
/usr/share/doc/apache/INSTALL.gz
/usr/share/doc/apache/KEYS.gz
/usr/share/doc/apache/config.layout.gz
/usr/share/doc/apache/README
/usr/share/doc/apache/README-WIN.TXT
/usr/share/doc/apache/WARNING-WIN.TXT
/usr/share/doc/apache/changelog.gz
/usr/share/doc/apache/intro.html
/usr/share/doc/apache/buildinfo.Debian
/usr/share/doc/apache/README.Debian
/usr/share/doc/apache/TODO.Debian
/usr/share/doc/apache/copyright
/usr/share/doc/apache/examples
/usr/share/doc/apache/examples/access.conf-dist
/usr/share/doc/apache/examples/highperformance.conf-dist
/usr/share/doc/apache/examples/httpd.conf-dist
/usr/share/doc/apache/examples/srm.conf-dist
/usr/share/doc/apache/examples/cron.conf
/usr/share/doc/apache/examples/httpd.conf
/usr/share/doc/apache/examples/access.conf
/usr/share/doc/apache/examples/srm.conf
/usr/share/doc/apache/examples/htaccess.debian
/usr/share/doc/apache/examples/magic.gz
/usr/share/doc/apache/examples/mime.types.gz
/usr/share/doc/apache/ABOUT_APACHE.gz
/usr/share/doc/apache/Announcement.gz
/usr/share/doc/apache/changelog.Debian.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/apache.8.gz
/usr/share/man/man8/apacheconfig.8.gz
/usr/share/man/man8/apachectl.8.gz
/var
/var/log
/var/log/apache
/var/cache
/var/cache/apache
