# var/lib/dpkg/info/apache.conffiles
/etc/init.d/apache
/etc/cron.daily/apache
