# var/lib/dpkg/info/bsdmainutils.list
/.
/usr
/usr/bin
/usr/bin/ncal
/usr/bin/col
/usr/bin/colcrt
/usr/bin/colrm
/usr/bin/column
/usr/bin/from
/usr/bin/hexdump
/usr/bin/look
/usr/bin/lorder
/usr/bin/ul
/usr/bin/write
/usr/bin/calendar
/usr/games
/usr/games/banner
/usr/share
/usr/share/calendar
/usr/share/calendar/calendar.all
/usr/share/calendar/calendar.birthday
/usr/share/calendar/calendar.christian
/usr/share/calendar/calendar.computer
/usr/share/calendar/calendar.croatian
/usr/share/calendar/calendar.german
/usr/share/calendar/calendar.history
/usr/share/calendar/calendar.holiday
/usr/share/calendar/calendar.judaic
/usr/share/calendar/calendar.music
/usr/share/calendar/calendar.russian
/usr/share/calendar/calendar.usholiday
/usr/share/calendar/calendar.world
/usr/share/calendar/calendar.pagan
/usr/share/calendar/calendar.discorian
/usr/share/calendar/de_DE.ISO_8859-1
/usr/share/calendar/de_DE.ISO_8859-1/calendar.all
/usr/share/calendar/de_DE.ISO_8859-1/calendar.feiertag
/usr/share/calendar/de_DE.ISO_8859-1/calendar.geschichte
/usr/share/calendar/de_DE.ISO_8859-1/calendar.kirche
/usr/share/calendar/de_DE.ISO_8859-1/calendar.literatur
/usr/share/calendar/de_DE.ISO_8859-1/calendar.musik
/usr/share/calendar/de_DE.ISO_8859-1/calendar.wissenschaft
/usr/share/calendar/hr_HR.ISO_8859-2
/usr/share/calendar/hr_HR.ISO_8859-2/calendar.all
/usr/share/calendar/hr_HR.ISO_8859-2/calendar.praznici
/usr/share/calendar/ru_RU.KOI8-R
/usr/share/calendar/ru_RU.KOI8-R/calendar.all
/usr/share/calendar/ru_RU.KOI8-R/calendar.common
/usr/share/calendar/ru_RU.KOI8-R/calendar.msk
/usr/share/calendar/ru_RU.KOI8-R/calendar.orthodox
/usr/share/calendar/ru_RU.KOI8-R/calendar.pagan
/usr/share/calendar/2001
/usr/share/calendar/2001/calendar.christian
/usr/share/calendar/2001/calendar.hindu
/usr/share/calendar/2001/calendar.judaic
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ul.1.gz
/usr/share/man/man1/look.1.gz
/usr/share/man/man1/hexdump.1.gz
/usr/share/man/man1/from.1.gz
/usr/share/man/man1/column.1.gz
/usr/share/man/man1/colrm.1.gz
/usr/share/man/man1/colcrt.1.gz
/usr/share/man/man1/col.1.gz
/usr/share/man/man1/write.1.gz
/usr/share/man/man1/lorder.1.gz
/usr/share/man/man1/ncal.1.gz
/usr/share/man/man1/calendar.1.gz
/usr/share/man/man3
/usr/share/man/man3/calendar.3.gz
/usr/share/man/man6
/usr/share/man/man6/banner.6.gz
/usr/share/doc
/usr/share/doc/bsdmainutils
/usr/share/doc/bsdmainutils/copyright
/usr/share/doc/bsdmainutils/changelog.Debian.gz
/etc
/etc/cron.daily
/etc/cron.daily/calendar
/etc/calendar
/etc/calendar/default
/usr/bin/cal
/usr/share/man/man1/cal.1.gz
