# var/lib/dpkg/info/bsdmainutils.conffiles
/etc/cron.daily/calendar
/etc/calendar/default
