# var/lib/dpkg/info/console-tools.list
/.
/bin
/bin/loadkeys
/sbin
/sbin/kbdrate
/usr
/usr/sbin
/usr/sbin/vcstime
/usr/sbin/writevt
/usr/sbin/setvesablank
/usr/bin
/usr/bin/psfaddtable
/usr/bin/psfgettable
/usr/bin/psfstriptable
/usr/bin/chvt
/usr/bin/deallocvt
/usr/bin/fgconsole
/usr/bin/screendump
/usr/bin/vt-is-UTF8
/usr/bin/openvt
/usr/bin/dumpkeys
/usr/bin/showkey
/usr/bin/setleds
/usr/bin/setmetamode
/usr/bin/kbd_mode
/usr/bin/getkeycodes
/usr/bin/setkeycodes
/usr/bin/mk_modmap
/usr/bin/consolechars
/usr/bin/showcfont
/usr/bin/unicode_start
/usr/bin/unicode_stop
/usr/bin/charset
/usr/bin/codepage
/usr/bin/splitfont
/usr/bin/fix_bs_and_del
/usr/bin/font2psf
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/chvt.1.gz
/usr/share/man/man1/deallocvt.1.gz
/usr/share/man/man1/fgconsole.1.gz
/usr/share/man/man1/kbd_mode.1.gz
/usr/share/man/man1/psfaddtable.1.gz
/usr/share/man/man1/psfgettable.1.gz
/usr/share/man/man1/setleds.1.gz
/usr/share/man/man1/setmetamode.1.gz
/usr/share/man/man1/showkey.1.gz
/usr/share/man/man1/unicode_start.1.gz
/usr/share/man/man1/unicode_stop.1.gz
/usr/share/man/man1/vt-is-UTF8.1.gz
/usr/share/man/man1/showcfont.1.gz
/usr/share/man/man1/charset.1.gz
/usr/share/man/man1/openvt.1.gz
/usr/share/man/man1/dumpkeys.1.gz
/usr/share/man/man1/loadkeys.1.gz
/usr/share/man/man1/psfstriptable.1.gz
/usr/share/man/man4
/usr/share/man/man4/console.4.gz
/usr/share/man/man4/console_ioctls.4.gz
/usr/share/man/man5
/usr/share/man/man5/keymaps.5.gz
/usr/share/man/man8
/usr/share/man/man8/getkeycodes.8.gz
/usr/share/man/man8/setkeycodes.8.gz
/usr/share/man/man8/kbdrate.8.gz
/usr/share/man/man8/consolechars.8.gz
/usr/share/doc
/usr/share/doc/console-tools-libs
/usr/share/doc/console-tools-libs/html
/usr/share/doc/console-tools-libs/html/lct-1.html
/usr/share/doc/console-tools-libs/html/lct-2.html
/usr/share/doc/console-tools-libs/html/lct-3.html
/usr/share/doc/console-tools-libs/html/lct-4.html
/usr/share/doc/console-tools-libs/html/lct-5.html
/usr/share/doc/console-tools-libs/html/lct-6.html
/usr/share/doc/console-tools-libs/html/lct-7.html
/usr/share/doc/console-tools-libs/html/lct-8.html
/usr/share/doc/console-tools-libs/html/lct.html
/usr/share/doc/console-tools-libs/examples
/usr/share/doc/console-tools-libs/examples/unicode
/usr/share/doc/console-tools-libs/examples/unicode/README.strange-name
/usr/share/doc/console-tools-libs/examples/unicode/♪♬
/usr/share/doc/console-tools-libs/examples/unicode/utflist.long.gz
/usr/share/doc/console-tools-libs/examples/unicode/utflist.cp437.gz
/usr/share/doc/console-tools-libs/examples/unicode/utflist.short.gz
/usr/share/doc/console-tools-libs/examples/screenmaps
/usr/share/doc/console-tools-libs/examples/screenmaps/latin1-16bit-hex
/usr/share/doc/console-tools-libs/examples/screenmaps/latin1-16bit-utf
/usr/share/doc/console-tools-libs/examples/sample.rc
/usr/share/doc/console-tools-libs/examples/spawn_console.c
/usr/share/doc/console-tools-libs/examples/spawn_login.c
/usr/share/doc/console-tools-libs/lct.sgml.gz
/usr/share/doc/console-tools-libs/lct.txt.gz
/usr/share/doc/console-tools-libs/contrib
/usr/share/doc/console-tools-libs/contrib/cirrus.videomodes
/usr/share/doc/console-tools-libs/contrib/keysyms.h.txt
/usr/share/doc/console-tools-libs/contrib/README.Crosser
/usr/share/doc/console-tools-libs/contrib/README.Euro
/usr/share/doc/console-tools-libs/contrib/set_kbd_repeat
/usr/share/doc/console-tools-libs/contrib/as400.kbd.gz
/usr/share/doc/console-tools-libs/contrib/console.docs.gz
/usr/share/doc/console-tools-libs/contrib/LISEZMOI.Euro.gz
/usr/share/doc/console-tools-libs/file-formats
/usr/share/doc/console-tools-libs/file-formats/cp
/usr/share/doc/console-tools-libs/file-formats/cpi
/usr/share/doc/console-tools-libs/file-formats/psf
/usr/share/doc/console-tools-libs/file-formats/raw
/usr/share/doc/console-tools-libs/file-formats/xpsf.draft.gz
/usr/share/doc/console-tools-libs/README.Debian
/usr/share/doc-base
/usr/share/doc-base/lct
/usr/share/keymaps
/usr/share/consolefonts
/usr/share/consoletrans
/usr/share/videomodes
/lib
/etc
/etc/console-tools
/etc/console-tools/config
/etc/init.d
/etc/init.d/console-screen.sh
/usr/bin/open
/usr/share/man/man1/open.1.gz
/usr/share/man/man1/codepage.1.gz
/usr/share/man/man1/fix_bs_and_del.1.gz
/usr/share/man/man1/mk_modmap.1.gz
/usr/share/man/man1/screendump.1.gz
/usr/share/man/man1/font2psf.1.gz
/usr/share/man/man1/splitfont.1.gz
/usr/share/man/man8/setvesablank.8.gz
/usr/share/man/man8/vcstime.8.gz
/usr/share/man/man8/writevt.8.gz
/usr/share/doc/console-tools
