# var/lib/dpkg/info/console-tools.conffiles
/etc/console-tools/config
/etc/init.d/console-screen.sh
