# var/lib/dpkg/info/cvs.conffiles
/etc/cron.weekly/cvs
/etc/cvs-cron.conf
/etc/cvs-pserver.conf
