# var/lib/dpkg/info/cvs.list
/.
/etc
/etc/cvs-cron.conf
/etc/cvs-pserver.conf
/etc/cron.weekly
/etc/cron.weekly/cvs
/usr
/usr/sbin
/usr/sbin/cvsconfig
/usr/sbin/cvs-makerepos
/usr/sbin/cvs-pserver
/usr/share
/usr/share/info
/usr/share/info/cvs.info-2.gz
/usr/share/info/cvs.info-3.gz
/usr/share/info/cvs.info-4.gz
/usr/share/info/cvs.info-5.gz
/usr/share/info/cvs.info-7.gz
/usr/share/info/cvs.info-8.gz
/usr/share/info/cvs.info-9.gz
/usr/share/info/cvsclient.info.gz
/usr/share/info/cvsclient.info-2.gz
/usr/share/info/cvsclient.info-3.gz
/usr/share/info/cvs.info.gz
/usr/share/info/cvs.info-1.gz
/usr/share/info/cvs.info-6.gz
/usr/share/info/cvsclient.info-1.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/cvsconfig.1.gz
/usr/share/man/man1/cvs.1.gz
/usr/share/man/man1/cvs-makerepos.1.gz
/usr/share/man/man5
/usr/share/man/man5/cvs.5.gz
/usr/share/man/man8
/usr/share/man/man8/cvsbug.8.gz
/usr/share/doc
/usr/share/doc/cvs
/usr/share/doc/cvs/README.gz
/usr/share/doc/cvs/changelog.gz
/usr/share/doc/cvs/README.Debian
/usr/share/doc/cvs/copyright
/usr/share/doc/cvs/ChangeLog.lib.gz
/usr/share/doc/cvs/ChangeLog.zlib.gz
/usr/share/doc/cvs/ChangeLog.fsf.lib
/usr/share/doc/cvs/ChangeLog.src.gz
/usr/share/doc/cvs/ChangeLog-97.src.gz
/usr/share/doc/cvs/ChangeLog-96.src.gz
/usr/share/doc/cvs/ChangeLog-9395.src.gz
/usr/share/doc/cvs/ChangeLog-9194.src.gz
/usr/share/doc/cvs/ChangeLog.man.gz
/usr/share/doc/cvs/ChangeLog.doc.gz
/usr/share/doc/cvs/ChangeLog.contrib.gz
/usr/share/doc/cvs/ChangeLog.fsf.doc
/usr/share/doc/cvs/ChangeLog.os2.gz
/usr/share/doc/cvs/ChangeLog.tools
/usr/share/doc/cvs/ChangeLog.emx.gz
/usr/share/doc/cvs/ChangeLog.vms.gz
/usr/share/doc/cvs/ChangeLog.zoo.gz
/usr/share/doc/cvs/FAQ.gz
/usr/share/doc/cvs/NEWS.gz
/usr/share/doc/cvs/changelog.Debian.gz
/usr/share/doc/cvs/ChangeLog.diff.gz
/usr/share/doc/cvs/ChangeLog.windows-NT.gz
/usr/share/doc/cvs/ChangeLog.windows-NT.SCC.gz
/usr/bin
/usr/bin/cvs
/usr/bin/cvsbug
/usr/bin/rcs2log
/usr/lib
/usr/lib/cvs
/usr/lib/cvs/contrib
/usr/lib/cvs/contrib/README
/usr/lib/cvs/contrib/intro.doc
/usr/lib/cvs/contrib/cvscheck.man
/usr/lib/cvs/contrib/clmerge
/usr/lib/cvs/contrib/cln_hist
/usr/lib/cvs/contrib/commit_prep
/usr/lib/cvs/contrib/cvs2vendor
/usr/lib/cvs/contrib/cvs_acls
/usr/lib/cvs/contrib/cvscheck
/usr/lib/cvs/contrib/log
/usr/lib/cvs/contrib/log_accum
/usr/lib/cvs/contrib/mfpipe
/usr/lib/cvs/contrib/rcs-to-cvs
/usr/lib/cvs/contrib/rcs2log
/usr/lib/cvs/contrib/rcslock
/usr/lib/cvs/contrib/sccs2rcs
