# var/lib/dpkg/info/devfsd.list
/.
/sbin
/sbin/devfsd
/etc
/etc/init.d
/etc/init.d/devfsd
/etc/devfs
/etc/devfs/symlinks.d
/etc/devfs/symlinks.d/devfsd
/etc/devfs/devices.d
/etc/devfs/devices.d/devfsd
/etc/devfs/conf.d
/etc/devfs/devfsd.conf
/etc/devfs/perms
/etc/devfs/symlinks
/etc/devfs/devices
/etc/devfs/dir.conf
/etc/modutils
/etc/modutils/devfsd
/usr
/usr/sbin
/usr/sbin/update-devfsd
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/devfsd.8.gz
/usr/share/doc
/usr/share/doc/devfsd
/usr/share/doc/devfsd/README.Debian
/usr/share/doc/devfsd/copyright
/usr/share/doc/devfsd/changelog.Debian.gz
