# var/lib/dpkg/info/devfsd.conffiles
/etc/devfs/devfsd.conf
/etc/devfs/symlinks
/etc/init.d/devfsd
/etc/devfs/perms
/etc/devfs/devices
/etc/devfs/devices.d/devfsd
/etc/devfs/symlinks.d/devfsd
/etc/devfs/dir.conf
