# var/lib/dpkg/info/findutils.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/findutils
/usr/share/doc/findutils/copyright
/usr/share/doc/findutils/changelog.Debian.gz
/usr/share/doc/findutils/ChangeLog.gz
/usr/share/doc/findutils/INSTALL.gz
/usr/share/doc/findutils/NEWS.gz
/usr/share/doc/findutils/README.gz
/usr/share/doc/findutils/TODO.gz
/usr/share/doc/findutils/README.debian
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/xargs.1.gz
/usr/share/man/man1/find.1.gz
/usr/share/man/man1/locate.1.gz
/usr/share/man/man1/updatedb.1.gz
/usr/share/man/man5
/usr/share/man/man5/locatedb.5.gz
/usr/share/info
/usr/share/info/find.info-2.gz
/usr/share/info/find.info.gz
/usr/share/info/find.info-1.gz
/usr/share/locale
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/findutils.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/findutils.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/findutils.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/findutils.mo
/usr/share/locale/ko
/usr/share/locale/ko/LC_MESSAGES
/usr/share/locale/ko/LC_MESSAGES/findutils.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/findutils.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/findutils.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/findutils.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/findutils.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/findutils.mo
/usr/share/locale/gl
/usr/share/locale/gl/LC_MESSAGES
/usr/share/locale/gl/LC_MESSAGES/findutils.mo
/usr/share/locale/et
/usr/share/locale/et/LC_MESSAGES
/usr/share/locale/et/LC_MESSAGES/findutils.mo
/usr/share/doc-base
/usr/share/doc-base/findutils
/usr/bin
/usr/bin/find
/usr/bin/xargs
/usr/bin/locate
/usr/bin/updatedb
/usr/lib
/usr/lib/locate
/usr/lib/locate/frcode
/usr/lib/locate/code
/usr/lib/locate/bigram
/etc
/etc/cron.daily
/etc/cron.daily/find
/etc/updatedb.conf
/var
/var/lib
/var/lib/locate
/usr/share/doc/findutils/changelog.gz
