# var/lib/dpkg/info/findutils.conffiles
/etc/cron.daily/find
/etc/updatedb.conf
