# var/lib/dpkg/info/gpm.list
/.
/usr
/usr/sbin
/usr/sbin/gpmconfig
/usr/sbin/gpm
/usr/sbin/gpm-mouse-test
/usr/sbin/gpm-microtouch-setup
/usr/bin
/usr/bin/gpm_has_mouse_control
/usr/bin/mev
/usr/bin/gpm-root
/usr/share
/usr/share/info
/usr/share/info/gpm.info.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mev.1.gz
/usr/share/man/man1/gpm-root.1.gz
/usr/share/man/man8
/usr/share/man/man8/gpm-mouse-test.8.gz
/usr/share/man/man8/gpm.8.gz
/usr/share/man/man8/gpmconfig.8.gz
/usr/share/doc
/usr/share/doc/gpm
/usr/share/doc/gpm/README.gunze.gz
/usr/share/doc/gpm/changelog.gz
/usr/share/doc/gpm/README.microtouch
/usr/share/doc/gpm/changelog.Debian.gz
/usr/share/doc/gpm/copyright
/usr/share/doc/gpm/faq.gz
/usr/share/doc/gpm/README.gz
/usr/share/doc/gpm/README.synaptics.gz
/usr/share/doc/gpm/README.twiddler.gz
/etc
/etc/init.d
/etc/init.d/gpm
/etc/gpm-root.conf
/usr/share/man/man1/gpm_has_mouse_control.1.gz
/usr/share/man/man8/gpm-microtouch-setup.8.gz
