# var/lib/dpkg/info/ifupdown.list
/.
/usr
/usr/bin
/usr/sbin
/usr/share
/usr/share/doc
/usr/share/doc/ifupdown
/usr/share/doc/ifupdown/examples
/usr/share/doc/ifupdown/examples/check-mac-address.sh
/usr/share/doc/ifupdown/examples/get-mac-address.sh
/usr/share/doc/ifupdown/examples/pcmcia-compat.sh
/usr/share/doc/ifupdown/examples/ping-places.sh
/usr/share/doc/ifupdown/examples/network-interfaces.gz
/usr/share/doc/ifupdown/copyright
/usr/share/doc/ifupdown/changelog.gz
/usr/share/doc/ifupdown/changelog.Debian.gz
/usr/share/ifupdown
/usr/share/ifupdown/upgrade-from-0.5.x.pl
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/interfaces.5.gz
/usr/share/man/man8
/usr/share/man/man8/ifup.8.gz
/etc
/etc/network
/etc/init.d
/etc/init.d/ifupdown
/sbin
/sbin/ifup
/sbin/ifdown
/usr/share/man/man8/ifdown.8.gz
