# var/lib/dpkg/info/ipchains.conffiles
/etc/default/ipchains
/etc/init.d/ipchains
