# var/lib/dpkg/info/ipchains.list
/.
/sbin
/sbin/ipchains
/sbin/ipchains-restore
/sbin/ipchains-save
/sbin/ipfwadm-wrapper
/usr
/usr/share
/usr/share/doc
/usr/share/doc/ipchains
/usr/share/doc/ipchains/changelog.gz
/usr/share/doc/ipchains/copyright
/usr/share/doc/ipchains/quickref.ps.gz
/usr/share/doc/ipchains/changelog.Debian.gz
/usr/share/man
/usr/share/man/man4
/usr/share/man/man4/ipfw_chains.4.gz
/usr/share/man/man4/ipfw-ipchainsalternative.4.gz
/usr/share/man/man8
/usr/share/man/man8/ipchains-save.8.gz
/usr/share/man/man8/ipfwadm-wrapper.8.gz
/usr/share/man/man8/ipfwadm-ipchainsalternative.8.gz
/usr/share/man/man8/ipchains.8.gz
/usr/share/man/man8/ipchains-restore.8.gz
/etc
/etc/default
/etc/default/ipchains
/etc/init.d
/etc/init.d/ipchains
