# var/lib/dpkg/info/klogd.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/klogd
/usr/share/doc/klogd/changelog.gz
/usr/share/doc/klogd/copyright
/usr/share/doc/klogd/modutils.patch
/usr/share/doc/klogd/changelog.Debian.gz
/usr/share/doc/klogd/readme.txt.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/klogd.8.gz
/sbin
/sbin/klogd
/etc
/etc/init.d
/etc/init.d/klogd
