# var/lib/dpkg/info/login.list
/.
/bin
/bin/login
/bin/su
/etc
/etc/init.d
/etc/init.d/logoutd
/etc/pam.d
/etc/pam.d/login
/etc/pam.d/su
/etc/login.defs
/etc/securetty
/usr
/usr/bin
/usr/bin/newgrp
/usr/bin/faillog
/usr/bin/lastlog
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/su.1.gz
/usr/share/man/man1/login.1.gz
/usr/share/man/man1/newgrp.1.gz
/usr/share/man/man5
/usr/share/man/man5/faillog.5.gz
/usr/share/man/man5/dialups.5.gz
/usr/share/man/man5/login.defs.5.gz
/usr/share/man/man8
/usr/share/man/man8/lastlog.8.gz
/usr/share/man/man8/logoutd.8.gz
/usr/share/man/man8/faillog.8.gz
/usr/share/doc
/usr/share/doc/login
/usr/share/doc/login/README.gz
/usr/share/doc/login/README.shadow-paper.gz
/usr/share/doc/login/README.pam.gz
/usr/share/doc/login/changelog.Debian.gz
/usr/share/doc/login/copyright
/usr/share/doc/login/changelog.gz
/usr/share/doc/login/README.debian.gz
/usr/sbin
/usr/sbin/logoutd
/usr/bin/sg
/usr/share/man/man1/sg.1.gz
/usr/share/man/man5/d_passwd.5.gz
