# var/lib/dpkg/info/login.conffiles
/etc/login.defs
/etc/pam.d/login
/etc/pam.d/su
/etc/init.d/logoutd
/etc/securetty
