# var/lib/dpkg/info/logrotate.conffiles
/etc/logrotate.conf
/etc/logrotate.conf
/etc/cron.daily/logrotate
