# var/lib/dpkg/info/modutils.list
/.
/sbin
/sbin/insmod
/sbin/modinfo
/sbin/insmod_ksymoops_clean
/sbin/kernelversion
/sbin/genksyms
/sbin/depmod
/sbin/update-modules
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/kernelversion.1.gz
/usr/share/man/man2
/usr/share/man/man2/get_kernel_syms.2.gz
/usr/share/man/man2/init_module.2.gz
/usr/share/man/man2/query_module.2.gz
/usr/share/man/man2/create_module.2.gz
/usr/share/man/man2/delete_module.2.gz
/usr/share/man/man5
/usr/share/man/man5/modules.conf.5.gz
/usr/share/man/man5/modules.5.gz
/usr/share/man/man8
/usr/share/man/man8/insmod.8.gz
/usr/share/man/man8/kallsyms.8.gz
/usr/share/man/man8/ksyms.8.gz
/usr/share/man/man8/lsmod.8.gz
/usr/share/man/man8/modinfo.8.gz
/usr/share/man/man8/rmmod.8.gz
/usr/share/man/man8/depmod.8.gz
/usr/share/man/man8/genksyms.8.gz
/usr/share/man/man8/modprobe.8.gz
/usr/share/man/man8/update-modules.8.gz
/usr/share/doc
/usr/share/doc/modutils
/usr/share/doc/modutils/NEWS.gz
/usr/share/doc/modutils/TODO.gz
/usr/share/doc/modutils/changelog.gz
/usr/share/doc/modutils/examples
/usr/share/doc/modutils/examples/versioned
/usr/share/doc/modutils/CREDITS.gz
/usr/share/doc/modutils/README.gz
/usr/share/doc/modutils/copyright
/usr/share/doc/modutils/module-policy.Debian.gz
/usr/share/doc/modutils/README.Debian.gz
/usr/share/doc/modutils/changelog.Debian.gz
/usr/bin
/lib
/lib/modules
/var
/var/log
/var/log/ksymoops
/etc
/etc/cron.daily
/etc/cron.daily/modutils
/etc/modules
/etc/modutils
/etc/modutils/0keep
/etc/modutils/actions
/etc/modutils/aliases
/etc/modutils/arch
/etc/modutils/arch/i386
/etc/modutils/arch/alpha
/etc/modutils/arch/m68k.generic
/etc/modutils/arch/m68k.mac
/etc/modutils/arch/m68k.amiga
/etc/modutils/arch/m68k.atari
/etc/modutils/paths
/etc/init.d
/etc/init.d/modutils
/sbin/rmmod
/sbin/modprobe
/sbin/lsmod
/sbin/ksyms
/sbin/kallsyms
/usr/bin/genksyms
