# var/lib/dpkg/info/netbase.conffiles
/etc/init.d/networking
/etc/gateways
/etc/protocols
/etc/services
/etc/rpc
