# var/lib/dpkg/info/netbase.list
/.
/usr
/usr/sbin
/usr/sbin/update-inetd
/usr/bin
/usr/share
/usr/share/doc
/usr/share/doc/netbase
/usr/share/doc/netbase/README.Debian
/usr/share/doc/netbase/copyright
/usr/share/doc/netbase/changelog.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man3
/usr/share/man/man3/DebianNet.3pm.gz
/usr/share/man/man4
/usr/share/man/man5
/usr/share/man/man5/gateways.5.gz
/usr/share/man/man8
/usr/share/man/man8/resolv+.8.gz
/usr/share/man/man8/update-inetd.8.gz
/usr/lib
/usr/lib/perl5
/usr/lib/perl5/Debian
/usr/lib/perl5/DebianNet.pm
/bin
/sbin
/etc
/etc/network
/etc/init.d
/etc/init.d/networking
/etc/cron.daily
/etc/gateways
/etc/protocols
/etc/rpc
/etc/services
