# var/lib/dpkg/info/netkit-inetd.list
/.
/etc
/etc/init.d
/etc/init.d/inetd
/etc/cron.daily
/etc/cron.daily/netkit-inetd
/usr
/usr/sbin
/usr/sbin/inetd
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/inetd.8.gz
/usr/share/doc
/usr/share/doc/netkit-inetd
/usr/share/doc/netkit-inetd/BUGS
/usr/share/doc/netkit-inetd/README
/usr/share/doc/netkit-inetd/copyright
/usr/share/doc/netkit-inetd/changelog.gz
/usr/share/doc/netkit-inetd/changelog.Debian.gz
