# var/lib/dpkg/info/netkit-inetd.conffiles
/etc/init.d/inetd
/etc/cron.daily/netkit-inetd
