# var/lib/dpkg/info/nfs-kernel-server.list
/.
/etc
/etc/init.d
/etc/init.d/nfs-kernel-server
/usr
/usr/sbin
/usr/sbin/exportfs
/usr/sbin/rpc.mountd
/usr/sbin/rpc.nfsd
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/exports.5.gz
/usr/share/man/man8
/usr/share/man/man8/mountd.8.gz
/usr/share/man/man8/exportfs.8.gz
/usr/share/man/man8/nfsd.8.gz
/usr/share/doc
/usr/share/doc/nfs-kernel-server
/usr/share/doc/nfs-kernel-server/README
/usr/share/doc/nfs-kernel-server/copyright
/usr/share/doc/nfs-kernel-server/ChangeLog.gz
/usr/share/doc/nfs-kernel-server/changelog.Debian.gz
/var
/var/lib
/var/lib/nfs
/usr/share/man/man8/rpc.mountd.8.gz
/usr/share/man/man8/rpc.nfsd.8.gz
