# var/lib/dpkg/info/nfs-common.list
/.
/etc
/etc/init.d
/etc/init.d/nfs-common
/sbin
/sbin/rpc.lockd
/sbin/rpc.statd
/sbin/showmount
/usr
/usr/sbin
/usr/sbin/nfsstat
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/showmount.8.gz
/usr/share/man/man8/lockd.8.gz
/usr/share/man/man8/statd.8.gz
/usr/share/man/man8/nfsstat.8.gz
/usr/share/doc
/usr/share/doc/nfs-common
/usr/share/doc/nfs-common/copyright
/usr/share/doc/nfs-common/changelog.Debian.gz
/var
/var/lib
/var/lib/nfs
/usr/share/man/man8/rpc.lockd.8.gz
/usr/share/man/man8/rpc.statd.8.gz
