# var/lib/dpkg/info/nvi-m17n.list
/.
/usr
/usr/bin
/usr/bin/nview-m17n
/usr/bin/nvi-m17n
/usr/bin/nex-m17n
/usr/lib
/usr/share
/usr/share/doc
/usr/share/doc/nvi-m17n
/usr/share/doc/nvi-m17n/docs.m17n
/usr/share/doc/nvi-m17n/docs.m17n/README.japanese.gz
/usr/share/doc/nvi-m17n/docs.m17n/changelog.gz
/usr/share/doc/nvi-m17n/docs.m17n/README.english.gz
/usr/share/doc/nvi-m17n/README.gz
/usr/share/doc/nvi-m17n/FAQ.gz
/usr/share/doc/nvi-m17n/changelog.Debian.gz
/usr/share/doc/nvi-m17n/LAYOUT
/usr/share/doc/nvi-m17n/README.Debian
/usr/share/doc/nvi-m17n/copyright
/usr/share/doc/nvi-m17n/vi.beginner.gz
/usr/share/doc/nvi-m17n/vi.advanced.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/nvi-m17n.1.gz
/etc
/etc/init.d
/etc/init.d/nvi-m17n
/usr/share/man/man1/nex-m17n.1.gz
/usr/share/man/man1/nview-m17n.1.gz
