# var/lib/dpkg/info/nvi-m17n.postrm
#! /bin/sh

if [ "$1" = "purge" ] ; then
	update-rc.d nvi-m17n remove >/dev/null
fi

   

