# var/lib/dpkg/info/nvi.list
/.
/usr
/usr/bin
/usr/bin/nvi
/usr/bin/nex
/usr/bin/nview
/usr/share
/usr/share/doc
/usr/share/doc/nvi
/usr/share/doc/nvi/README.gz
/usr/share/doc/nvi/FAQ.gz
/usr/share/doc/nvi/changelog.Debian.gz
/usr/share/doc/nvi/README.Debian
/usr/share/doc/nvi/copyright
/usr/share/doc/nvi/changelog.gz
/usr/share/doc/nvi/vi.beginner.gz
/usr/share/doc/nvi/vi.advanced.gz
/usr/share/vi
/usr/share/vi/catalog
/usr/share/vi/catalog/dutch
/usr/share/vi/catalog/english
/usr/share/vi/catalog/french
/usr/share/vi/catalog/german
/usr/share/vi/catalog/ru_SU.KOI8-R
/usr/share/vi/catalog/spanish
/usr/share/vi/catalog/swedish
/usr/share/vi/recover
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/nvi.1.gz
/etc
/etc/init.d
/etc/init.d/nviboot
/usr/share/man/man1/nex.1.gz
/usr/share/man/man1/nview.1.gz
