# var/lib/dpkg/info/portmap.list
/.
/bin
/sbin
/sbin/portmap
/sbin/pmap_dump
/sbin/pmap_set
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/pmap_set.8.gz
/usr/share/man/man8/portmap.8.gz
/usr/share/man/man8/pmap_dump.8.gz
/usr/share/doc
/usr/share/doc/portmap
/usr/share/doc/portmap/changelog.gz
/usr/share/doc/portmap/portmapper.txt.gz
/usr/share/doc/portmap/copyright
/usr/share/doc/portmap/README.gz
/usr/share/doc/portmap/changelog.Debian.gz
/etc
/etc/init.d
/etc/init.d/portmap
