# var/lib/dpkg/info/sysklogd.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/sysklogd
/usr/share/doc/sysklogd/changelog.gz
/usr/share/doc/sysklogd/copyright
/usr/share/doc/sysklogd/changelog.Debian.gz
/usr/share/doc/sysklogd/readme.txt.gz
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/syslog.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/syslogd.8.gz
/usr/share/man/man8/sysklogd.8.gz
/usr/share/man/man8/syslog-facility.8.gz
/usr/share/man/man8/syslogd-listfiles.8.gz
/usr/sbin
/usr/sbin/syslogd-listfiles
/usr/sbin/syslog-facility
/sbin
/sbin/syslogd
/etc
/etc/init.d
/etc/init.d/sysklogd
/etc/cron.daily
/etc/cron.daily/sysklogd
/etc/cron.weekly
/etc/cron.weekly/sysklogd
/etc/syslog.conf
/var
/var/log
