# var/lib/dpkg/info/sysklogd.conffiles
/etc/syslog.conf
/etc/init.d/sysklogd
/etc/cron.daily/sysklogd
/etc/cron.weekly/sysklogd
