# var/lib/dpkg/info/sysvinit.list
/.
/etc
/etc/default
/etc/init.d
/etc/init.d/halt
/etc/init.d/rc
/etc/init.d/reboot
/etc/init.d/single
/etc/init.d/skeleton
/etc/init.d/README
/etc/init.d/rmnologin
/etc/init.d/urandom
/etc/init.d/rcS
/etc/init.d/bootmisc.sh
/etc/init.d/checkroot.sh
/etc/init.d/checkfs.sh
/etc/init.d/mountall.sh
/etc/init.d/mountnfs.sh
/etc/init.d/hostname.sh
/etc/init.d/sendsigs
/etc/init.d/umountfs
/etc/init.d/umountnfs.sh
/etc/rc.boot
/etc/rc0.d
/etc/rc1.d
/etc/rc2.d
/etc/rc3.d
/etc/rc4.d
/etc/rc5.d
/etc/rc6.d
/etc/rcS.d
/etc/rcS.d/README
/usr
/usr/doc
/usr/share
/usr/share/doc
/usr/share/doc/sysvinit
/usr/share/doc/sysvinit/examples
/usr/share/doc/sysvinit/examples/README
/usr/share/doc/sysvinit/examples/bootlogd
/usr/share/doc/sysvinit/examples/default.rcS
/usr/share/doc/sysvinit/examples/network
/usr/share/doc/sysvinit/examples/powerfail
/usr/share/doc/sysvinit/examples/inittab
/usr/share/doc/sysvinit/copyright
/usr/share/doc/sysvinit/changelog.gz
/usr/share/doc/sysvinit/README.runlevels.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mesg.1.gz
/usr/share/man/man1/last.1.gz
/usr/share/man/man1/lastb.1.gz
/usr/share/man/man5
/usr/share/man/man5/inittab.5.gz
/usr/share/man/man5/rc.boot.5.gz
/usr/share/man/man5/rcS.5.gz
/usr/share/man/man5/initscript.5.gz
/usr/share/man/man8
/usr/share/man/man8/killall5.8.gz
/usr/share/man/man8/pidof.8.gz
/usr/share/man/man8/poweroff.8.gz
/usr/share/man/man8/reboot.8.gz
/usr/share/man/man8/shutdown.8.gz
/usr/share/man/man8/sulogin.8.gz
/usr/share/man/man8/telinit.8.gz
/usr/share/man/man8/halt.8.gz
/usr/share/man/man8/init.8.gz
/usr/share/man/man8/runlevel.8.gz
/usr/include
/usr/include/initreq.h
/usr/bin
/usr/bin/mesg
/usr/bin/last
/bin
/sbin
/sbin/fsck.nfs
/sbin/halt
/sbin/init
/sbin/killall5
/sbin/sulogin
/sbin/runlevel
/sbin/shutdown
/usr/bin/lastb
/bin/pidof
/sbin/reboot
/sbin/poweroff
/sbin/telinit
