# var/lib/dpkg/info/sysvinit.conffiles
/etc/init.d/bootmisc.sh
/etc/init.d/checkfs.sh
/etc/init.d/checkroot.sh
/etc/init.d/halt
/etc/init.d/hostname.sh
/etc/init.d/mountall.sh
/etc/init.d/mountnfs.sh
/etc/init.d/reboot
/etc/init.d/rmnologin
/etc/init.d/sendsigs
/etc/init.d/single
/etc/init.d/umountfs
/etc/init.d/urandom
