# var/lib/dpkg/info/util-linux.list
/.
/bin
/bin/arch
/bin/dmesg
/bin/more
/sbin
/sbin/mkswap
/sbin/hwclock
/sbin/fsck.minix
/sbin/mkfs
/sbin/mkfs.minix
/sbin/getty
/etc
/etc/init.d
/etc/init.d/hwclock.sh
/usr
/usr/share
/usr/share/locale
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/util-linux.mo
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/util-linux.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/util-linux.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/util-linux.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/util-linux.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/util-linux.mo
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/ja/LC_MESSAGES/util-linux.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/util-linux.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/util-linux.mo
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/rev.1.gz
/usr/share/man/man1/mcookie.1.gz
/usr/share/man/man1/arch.1.gz
/usr/share/man/man1/readprofile.1.gz
/usr/share/man/man1/chkdupexe.1.gz
/usr/share/man/man1/setterm.1.gz
/usr/share/man/man1/whereis.1.gz
/usr/share/man/man1/ddate.1.gz
/usr/share/man/man1/getopt.1.gz
/usr/share/man/man1/more.1.gz
/usr/share/man/man1/namei.1.gz
/usr/share/man/man8
/usr/share/man/man8/ipcs.8.gz
/usr/share/man/man8/setsid.8.gz
/usr/share/man/man8/getty.8.gz
/usr/share/man/man8/hwclock.8.gz
/usr/share/man/man8/mkswap.8.gz
/usr/share/man/man8/fsck.minix.8.gz
/usr/share/man/man8/mkfs.minix.8.gz
/usr/share/man/man8/mkfs.8.gz
/usr/share/man/man8/dmesg.8.gz
/usr/share/man/man8/ipcrm.8.gz
/usr/share/man/man8/elvtune.8.gz
/usr/share/info
/usr/share/info/ipc.info.gz
/usr/share/doc
/usr/share/doc/util-linux
/usr/share/doc/util-linux/examples
/usr/share/doc/util-linux/examples/test.bash.gz
/usr/share/doc/util-linux/examples/test.tcsh.gz
/usr/share/doc/util-linux/examples/parse.bash.gz
/usr/share/doc/util-linux/examples/parse.tcsh.gz
/usr/share/doc/util-linux/examples/fstab.example2.gz
/usr/share/doc/util-linux/examples/sfdisk.examples.gz
/usr/share/doc/util-linux/README.poeigl.gz
/usr/share/doc/util-linux/README.Debian.hwclock.gz
/usr/share/doc/util-linux/HISTORY.gz
/usr/share/doc/util-linux/README.modems-with-getty.gz
/usr/share/doc/util-linux/copyright
/usr/share/doc/util-linux/README.getty.gz
/usr/share/doc/util-linux/changelog.Debian.gz
/usr/lib
/usr/lib/mime
/usr/lib/mime/packages
/usr/lib/mime/packages/util-linux
/usr/lib/more.help
/usr/sbin
/usr/sbin/readprofile
/usr/sbin/elvtune
/usr/bin
/usr/bin/ipcs
/usr/bin/ipcrm
/usr/bin/setsid
/usr/bin/namei
/usr/bin/setterm
/usr/bin/mcookie
/usr/bin/whereis
/usr/bin/ddate
/usr/bin/getopt
/usr/bin/rev
/usr/bin/chkdupexe
/usr/bin/fdformat
/usr/share/man/man8/clock.8.gz
