# etc/cron.daily/modutils
#! /bin/sh

[ -e /sbin/insmod_ksymoops_clean ] || exit 0
/sbin/insmod_ksymoops_clean

